# cron.daily/run_apt_system_mx
#!/bin/bash

sudo apt-get clean
sudo apt-get autoclean
sudo apt-get autoremove -y
sudo update-grub

